# systemd/units/bumblebee/podman-prune.timer
[Unit]
Description=Weekly Podman image prune

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
