# systemd/units/optimusprime/amdgpu-metrics.service
[Unit]
Description=Collect AMD GPU metrics for Prometheus

[Service]
Type=oneshot
ExecStart=/usr/local/bin/amdgpu-metrics.sh
