# systemd/containers/system/bumblebee/notify-recovery-check.service
[Unit]
Description=Check for recovered services and notify via Telegram

[Service]
Type=oneshot
ExecStart=/home/matteo/notify-recovery-check.sh
