# systemd/containers/system/bumblebee/openclaw.container
[Unit]
Description=OpenClaw AI Agent
Requires=network-online.target
Requires=ollama.service
After=ollama.service

[Container]
ContainerName=openclaw
Image=ghcr.io/openclaw/openclaw:latest
PublishPort=18789:18789
Volume=/etc/localtime:/etc/localtime:ro
Volume=/home/matteo/docker_persistent/openclaw/config:/home/node/.openclaw:rw,Z
Volume=/home/matteo/docker_persistent/openclaw/workspace:/home/node/.openclaw/workspace:rw,Z
Environment="LLM_MODEL=ollama/qwen2.5-coder:7b"
Environment="LLM_BASE_URL=http://192.168.1.14:11434"
Environment="OLLAMA_BASE_URL=http://192.168.1.14:11434"
Environment="OPENCLAW_HOST=0.0.0.0"

AutoUpdate=registry

Label=traefik.enable=true
Label=traefik.http.routers.openclaw.rule=Host(`openclaw.bumblebee.favarohome.com`)
Label=traefik.http.routers.openclaw.entrypoints=websecure
Label=traefik.http.routers.openclaw.tls.certresolver=cloudflare
Label=traefik.http.services.openclaw.loadbalancer.server.port=18789

[Service]
Restart=always

[Install]
WantedBy=default.target
