# systemd/containers/system/bumblebee/opencode.container
[Unit]
OnFailure=notify-failure@%n.service
Description=OpenCode AI Coding Assistant
Requires=network-online.target
Requires=ollama.service
After=ollama.service

[Container]
ContainerName=opencode
Image=ghcr.io/joostme/opencode-docker:latest
PublishPort=4096:4096
PublishPort=8081:8080
Volume=/etc/localtime:/etc/localtime:ro
Volume=/home/matteo/docker_persistent/opencode/config:/home/coder/.config:rw,Z
Volume=/home/matteo/docker_persistent/opencode/share:/home/coder/.local/share:rw,Z
Volume=/home/matteo/docker_persistent/opencode/repos:/home/coder/repos:rw,Z
Environment="OLLAMA_HOST=http://192.168.1.14:11434"

AutoUpdate=registry

Label=traefik.enable=true
Label=traefik.http.routers.opencode.rule=Host(`opencode.bumblebee.favarohome.com`)
Label=traefik.http.routers.opencode.entrypoints=websecure
Label=traefik.http.routers.opencode.tls.certresolver=cloudflare
Label=traefik.http.services.opencode.loadbalancer.server.port=4096

[Service]
Restart=always

[Install]
WantedBy=default.target
