# systemd/containers/system/bumblebee/opencode.container
[Unit]
OnFailure=notify-failure@%n.service
OnSuccess=notify-success@%n.service
Description=OpenCode AI Coding Assistant
Requires=network-online.target
Requires=ollama.service
After=ollama.service

[Container]
ContainerName=opencode
Image=ghcr.io/joostme/opencode-docker:latest
PublishPort=4096:4096
PublishPort=8081:8080
Volume=/etc/localtime:/etc/localtime:ro
Volume=/home/matteo/docker_persistent/opencode/home:/home/opencode:rw,Z
Volume=/home/matteo/docker_persistent/opencode/repos:/repos:rw,Z
Volume=/home/matteo/docker_persistent/opencode/local-bin:/usr/local/bin:rw,Z
Environment="OPENAI_API_KEY=sk-b6a7de9570a144ba89649f59f2222c8d"
Environment="OPENAI_BASE_URL=http://192.168.1.14:3000/api"

AutoUpdate=registry

Label=traefik.enable=true
Label=traefik.http.routers.opencode.rule=Host(`opencode.bumblebee.favarohome.com`)
Label=traefik.http.routers.opencode.entrypoints=websecure
Label=traefik.http.routers.opencode.tls.certresolver=cloudflare
Label=traefik.http.services.opencode.loadbalancer.server.port=4096

[Service]
Restart=always

[Install]
WantedBy=default.target
