# systemd/containers/system/bumblebee/openclaw.container
[Unit]
OnFailure=notify-failure@%n.service
OnSuccess=notify-success@%n.service
Description=OpenClaw AI Agent
Requires=network-online.target
Wants=ollama.service
After=ollama.service

[Container]
ContainerName=openclaw
Image=ghcr.io/openclaw/openclaw:latest
PublishPort=18789:18789
SecurityLabelDisable=true
PodmanArgs=--group-add 0
Volume=/etc/localtime:/etc/localtime:ro
Volume=/run/podman/podman.sock:/var/run/docker.sock
Volume=/home/matteo/openclaw-workspace:/home/node/.openclaw/workspace:rw
Volume=/home/matteo/docker_persistent/openclaw/config:/home/node/.openclaw:rw
Volume=/home/matteo/docker_persistent/openclaw/bin:/opt/agents:rw
Volume=/home/matteo/docker_persistent/openclaw/opencode-config:/home/node/.config/opencode:rw
Volume=/home/matteo/docker_persistent/openclaw/ssh:/home/node/.ssh:rw
Environment="OPENCLAW_HOST=0.0.0.0"
Environment="PATH=/opt/agents:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PLAYWRIGHT_BROWSERS_PATH=/opt/agents/playwright"
Environment="TZ=Europe/Zurich"
Environment="GIT_CONFIG_GLOBAL=/home/node/.openclaw/.gitconfig"
EnvironmentFile=/etc/containers/secrets/openclaw.env

AutoUpdate=registry


[Service]
Restart=always
ExecStartPost=/bin/bash -c "sleep 5 && /usr/bin/podman exec --user root openclaw /bin/bash /opt/agents/setup-playwright.sh; /usr/bin/podman exec --user root openclaw /bin/bash /opt/agents/setup-opencode.sh; /usr/bin/podman exec --user root openclaw /bin/bash /opt/agents/patch-models.sh; chattr +i /home/matteo/docker_persistent/openclaw/config/agents/main/agent/models.json 2>/dev/null || true"

[Install]
WantedBy=default.target
