# systemd/containers/system/bumblebee/paperless.container
[Unit]
OnFailure=notify-failure@%n.service
Description=Paperless-NGX document management
Requires=network-online.target
After=paperless-db.service paperless-broker.service paperless-gotenberg.service paperless-tika.service

[Container]
ContainerName=paperless
Network=paperless.network
Network=podman
Image=ghcr.io/paperless-ngx/paperless-ngx:latest
PublishPort=8000:8000
Volume=/etc/localtime:/etc/localtime:ro
Volume=/home/matteo/docker_persistent/paperless/data:/usr/src/paperless/data:rw,Z
Volume=/home/matteo/docker_persistent/paperless/media:/usr/src/paperless/media:rw,Z
Volume=/home/matteo/docker_persistent/paperless/export:/usr/src/paperless/export:rw,Z
Volume=/home/matteo/docker_persistent/paperless/consume:/usr/src/paperless/consume:rw,Z

Environment=TZ=Europe/Zurich
Environment=PAPERLESS_TIME_ZONE=Europe/Zurich
Environment=PAPERLESS_URL=https://paperless.bumblebee.favarohome.com
Environment=PAPERLESS_REDIS=redis://paperless-broker:6379
Environment=PAPERLESS_DBHOST=paperless-db
Environment=PAPERLESS_DBENGINE=postgresql
Environment=PAPERLESS_DBNAME=paperless
Environment=PAPERLESS_DBUSER=paperless
Environment=PAPERLESS_TIKA_ENABLED=1
Environment=PAPERLESS_TIKA_GOTENBERG_ENDPOINT=http://paperless-gotenberg:3000
Environment=PAPERLESS_TIKA_ENDPOINT=http://paperless-tika:9998
Environment=PAPERLESS_OCR_LANGUAGE=eng+ita
Environment=USERMAP_UID=1000
Environment=USERMAP_GID=1000

EnvironmentFile=/etc/containers/secrets/paperless.env

HealthCmd=CMD-SHELL curl -sf http://localhost:8000/ || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=120s
HealthTimeout=10s
HealthOnFailure=restart

AutoUpdate=registry

Label=traefik.enable=true
Label=traefik.http.routers.paperless.rule=Host(`paperless.bumblebee.favarohome.com`)
Label=traefik.http.routers.paperless.entrypoints=websecure
Label=traefik.http.routers.paperless.tls.certresolver=cloudflare
Label=traefik.http.services.paperless.loadbalancer.server.port=8000

[Service]
Restart=always

[Install]
WantedBy=default.target
