# systemd/containers/system/bumblebee/podman-auto-update-notify.timer
[Unit]
Description=Podman auto-update with notification timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
