# systemd/containers/system/bumblebee/podman-auto-update-notify.service
[Unit]
Description=Podman auto-update with Telegram notification
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
  output=$(/usr/bin/podman auto-update 2>&1); \
  rc=$?; \
  updated=$(echo "$output" | grep -c "true" || true); \
  if [ "$updated" -gt 0 ]; then \
    /usr/local/bin/telegram-send --config /etc/telegram-send.conf "[bumblebee] podman auto-update completed, $updated container(s) updated:\n$output"; \
  else \
    /usr/local/bin/telegram-send --config /etc/telegram-send.conf "[bumblebee] podman auto-update completed, no containers updated"; \
  fi; \
  exit $rc'
