# systemd/containers/system/optimusprime/grocy.container
[Unit]
OnFailure=notify-failure@%n.service
Description=Grocy - grocery & household management
Requires=network-online.target

[Container]
ContainerName=grocy
Image=lscr.io/linuxserver/grocy:latest
PublishPort=9283:80
Volume=/etc/localtime:/etc/localtime:ro
Volume=/mnt/data/docker_persistent/grocy/config:/config:rw

Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=Europe/Zurich

AutoUpdate=registry

HealthCmd=CMD-SHELL curl -sf -o /dev/null http://localhost:80/ || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=60s
HealthTimeout=10s
HealthOnFailure=restart

Label=traefik.enable=true
Label=traefik.http.routers.grocy.rule=Host(`grocy.optimusprime.favarohome.com`)
Label=traefik.http.routers.grocy.entrypoints=websecure
Label=traefik.http.routers.grocy.tls.certresolver=cloudflare
Label=traefik.http.services.grocy.loadbalancer.server.port=80

[Service]
Restart=always

[Install]
WantedBy=default.target
