# systemd/containers/system/optimusprime/influxdb.container
[Unit]
OnFailure=notify-failure@%n.service
Description=InfluxDB Time-Series Database
Requires=network-online.target

[Container]
ContainerName=influxdb
Image=docker.io/influxdb:2
PublishPort=8086:8086
Volume=/etc/localtime:/etc/localtime:ro
Volume=/mnt/data/docker_persistent/influxdb/data:/var/lib/influxdb2:rw
Volume=/mnt/data/docker_persistent/influxdb/config:/etc/influxdb2:rw
Environment=DOCKER_INFLUXDB_INIT_MODE=setup
Environment=DOCKER_INFLUXDB_INIT_USERNAME=admin
Environment=DOCKER_INFLUXDB_INIT_ORG=favarohome
Environment=DOCKER_INFLUXDB_INIT_BUCKET=homeassistant
Environment=DOCKER_INFLUXDB_INIT_RETENTION=0
EnvironmentFile=/etc/containers/secrets/influxdb.env

AutoUpdate=registry

HealthCmd=CMD-SHELL curl -sf -o /dev/null http://localhost:8086/health || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=60s
HealthTimeout=10s
HealthOnFailure=restart

Label=traefik.enable=true
Label=traefik.http.routers.influxdb.rule=Host(`influxdb.optimusprime.favarohome.com`)
Label=traefik.http.routers.influxdb.entrypoints=websecure
Label=traefik.http.routers.influxdb.tls.certresolver=cloudflare
Label=traefik.http.services.influxdb.loadbalancer.server.port=8086

[Service]
Restart=always

[Install]
WantedBy=default.target
