# systemd/containers/system/optimusprime/notify-failure@.service
[Unit]
Description=Send Telegram notification on service failure

[Service]
Type=oneshot
User=matteo
ExecStart=/home/matteo/notify-failure.sh %i
