# systemd/containers/system/optimusprime/notify-recovery-check.timer
[Unit]
Description=Periodically check for recovered services

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
