# systemd/containers/system/optimusprime/plex-rebuild.timer
[Unit]
Description=Daily check for upstream Plex updates and rebuild the local image

[Timer]
OnCalendar=*-*-* 04:45:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
