# systemd/containers/system/optimusprime/podman-auto-update-notify.service
[Unit]
Description=Podman auto-update with Telegram notification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/podman-update-notify.sh
