# systemd/containers/system/optimusprime/podman-auto-update-notify.service
[Unit]
Description=Podman auto-update with Telegram notification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'result=$(podman auto-update 2>&1); updated=$(echo "$result" | grep -c "true"); echo "$result"; if [ "$updated" -gt 0 ]; then telegram-send "[optimus-prime] podman auto-update: $updated container(s) updated."; else telegram-send "[optimus-prime] podman auto-update: all containers up to date."; fi'
