# systemd/containers/system/optimusprime/shutter-encoder.container
[Unit]
OnFailure=notify-failure@%n.service
Description=Shutter Encoder Container
Requires=network-online.target

[Container]
ContainerName=shutter-encoder
Image=docker.io/ciberkids/shutter-encoder:latest
PublishPort=6900:6080
Volume=/etc/localtime:/etc/localtime:ro
Volume=/mnt/data/docker_persistent/shutter-encoder:/config:rw
Volume=/mnt/MovieAndTvShows/ToFix:/input:ro
Volume=/mnt/MovieAndTvShows/ToFix/shutterOutput:/output:rw
AddDevice=/dev/dri/renderD128

Environment=TZ=Europe/Zurich
Environment=PUID=1000
Environment=PGID=1000
# GPU hw-encode drivers baked into the image:
#   AMD VA-API:   radeonsi (GCN/RDNA), r600 (TeraScale)  → mesa-va-drivers
#   Intel VA-API: iHD (Broadwell+), i965 (Haswell-)      → intel-media-va-driver, i965-va-driver
#   Intel QSV:    via libmfx1 shim                        → libmfx1
#   NVIDIA NVENC: runtime-injected by NVIDIA Container Toolkit (no image pkgs)
# Override LIBVA_DRIVER_NAME at runtime: radeonsi (AMD) | iHD (Intel new) | i965 (Intel old)
# For NVIDIA, add: AddDevice=nvidia.com/gpu=all (and remove /dev/dri if no iGPU)
Environment=LIBVA_DRIVER_NAME=radeonsi

AutoUpdate=registry

HealthCmd=CMD-SHELL curl -sf -o /dev/null http://localhost:6080/ || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=60s
HealthTimeout=10s
HealthOnFailure=restart

Label=traefik.enable=true
Label=traefik.http.routers.shutter-encoder.rule=Host(`shutter-encoder.optimusprime.favarohome.com`)
Label=traefik.http.routers.shutter-encoder.entrypoints=websecure
Label=traefik.http.routers.shutter-encoder.tls.certresolver=cloudflare
Label=traefik.http.services.shutter-encoder.loadbalancer.server.port=6080

[Service]
Restart=always

[Install]
WantedBy=default.target
