# systemd/containers/system/optimusprime/warracker.container
[Unit]
Description=Warracker warranty tracker
Requires=network-online.target
After=warracker-db.service

[Container]
ContainerName=warracker
Network=warracker.network
Image=ghcr.io/sassanix/warracker/main:latest
PublishPort=8005:80
Volume=/etc/localtime:/etc/localtime:ro
Volume=/mnt/data/docker_persistent/warracker/uploads:/data/uploads:rw

Environment=DB_HOST=warracker-db
Environment=DB_NAME=warranty_db
Environment=DB_USER=warranty_user
Environment=FRONTEND_URL=https://warracker.optimusprime.favarohome.com
Environment=APP_BASE_URL=https://warracker.optimusprime.favarohome.com
Environment=FLASK_ENV=production
Environment=WARRACKER_MEMORY_MODE=optimized
Environment=MAX_UPLOAD_MB=32
EnvironmentFile=/etc/containers/secrets/warracker.env

AutoUpdate=registry

Label=traefik.enable=true
Label=traefik.http.routers.warracker.rule=Host(`warracker.optimusprime.favarohome.com`)
Label=traefik.http.routers.warracker.entrypoints=websecure
Label=traefik.http.routers.warracker.tls.certresolver=cloudflare
Label=traefik.http.services.warracker.loadbalancer.server.port=80

[Service]
Restart=always

[Install]
WantedBy=default.target
