# systemd/services/nvidiapodmanchecks/checkNvidiaPodman.timer
[Unit]
Description=Run Nvidia Podman Check Every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target 
