# systemd/system/bumblebee/runner-cache-cleanup.service
[Unit]
Description=Clean GitLab runner cache entries older than 14 days

[Service]
Type=oneshot
ExecStart=/usr/local/bin/runner-cache-cleanup.sh
StandardOutput=journal
StandardError=journal
